# packaging/resources/ovirt_etl_watchdog.cron
#!/bin/bash
#


# ETL functions library.
. /etc/init.d/ovirt-etl

if [ ! -e /usr/share/ovirt-dwh/etl/kill ]; then
    pid=0
    procrunning > /dev/null 2>&1
    if [ $pid = '0' ]; then
        service ovirt-dwhd start > /dev/null 2>&1
    fi
fi
